# pi-gen/stage-flipperpi/files/rc.local
#!/bin/bash
python3 /home/pi/flipperpi/main.py &
exit 0
